# atom/.atom/recovery/celery_worker-668afd.conf
[program:celery_worker]
command=/home/ubuntu/.virtualenvs/cosmic-cloud/bin/celery worker -A adder -l info
directory=/home/ubuntu/cosmic-cloud
environment=PATH=/home/ubuntu/.virtualenvs/cosmic-cloud/bin:%(ENV_PATH)s,VIRTUAL_ENV="/home/ubuntu/.virtualenvs/cosmic-cloud"
user=ubuntu
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
killasgroup=true
startretries=1
stdout_logfile=/var/log/celery_worker/celery_worker.log
stderr_logfile=/var/log/celery_worker/celery_worker.log
