# atom/.atom/recovery/celery_worker-132456.conf
[program:celery_worker]
command=/home/ubuntu/.virtualenvs/cosmic-cloud/bin/celery worker -A adder -l info
directory=/home/ubuntu/cosmic-cloud
user=ubuntu
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
killasgroup=true
startretries=1
stdout_logfile=/home/ubuntu/logs/celery_worker.log
stderr_logfile=/home/ubuntu/logs/celery_worker.log
